# supervisord.conf
[supervisord]
nodaemon=true

[program:cron]
command = cron -f

[program:apache2]
command=/usr/sbin/apachectl -D FOREGROUND
autostart=true
autorestart=true