# docker/etc/php/supervisord.conf
[supervisord]
nodaemon=true

[program:cron]
command = cron -f

[program:php-fpm]
command = docker-php-entrypoint php-fpm
